# mapcrafter-config/mapcrafter.timer
[Unit]
Description=Update MapCrafter Minecraft map

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

