# ddns-tool/ddns.timer
[Unit]
Description=Update dynamic DNS

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

